# external/adhd/init/cras.conf
# Copyright (c) 2012 The Chromium OS Authors. All rights reserved.
# Use of this source code is governed by a BSD-style license that can be
# found in the LICENSE file.

# Installed by ADHD package.
# cras upstart job.

description     "ChromeOS audio server"
author          "chromium-os-dev@chromium.org"

env CRAS_SOCKET_DIR=/run/cras

start on starting system-services
stop on stopping system-services
respawn

# Allow the audio server real time priority.
limit rtprio 12 12

pre-start script
  mkdir -p -m 1770 "${CRAS_SOCKET_DIR}"
  chown -R cras:cras "${CRAS_SOCKET_DIR}"
end script

exec /bin/sh /usr/share/cros/init/cras.sh
